# marketplace-scanner-nz.service
[Unit]
Description=XCT NZ New Releases Scanner
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker run --rm --network host \
    --name xct-scanner-nz \
    -e DATABASE_URL=postgresql://freshdex:998d73bad57169c5b68a8a929e4171f6@localhost/freshdex_cdn \
    -e SCANNER_ACCOUNT_DIR=/app/scanner_account \
    -v /opt/xct-live/scanner_account:/app/scanner_account \
    freshdex-xct-live python3 marketplace_scanner.py scan-nz
ExecStop=/usr/bin/docker stop xct-scanner-nz
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
